# node-upstart.conf
start on started
stop on shutdown

exec node /home/ec2-user/ex2/index.js